# deploy/systemd/proethica-web.service
[Unit]
Description=ProEthica Web Application
After=network.target postgresql.service
Requires=postgresql.service
PartOf=proethica.target

[Service]
Type=exec
User=chris
Group=chris
WorkingDirectory=/opt/proethica
Environment=PATH=/opt/proethica/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=DATABASE_URL=postgresql://proethica_user:PROETHICA_PASS_HERE@localhost:5432/ai_ethical_dm
Environment=SQLALCHEMY_DATABASE_URI=postgresql://proethica_user:PROETHICA_PASS_HERE@localhost:5432/ai_ethical_dm
Environment=FLASK_ENV=production
Environment=SECRET_KEY=GENERATED_SECRET_HERE
Environment=USE_ONTSERVE=true
Environment=ONTSERVE_URL=http://localhost:8083
Environment=ONTSERVE_WEB_URL=http://localhost:5003
Environment=ENHANCED_SCENARIO_GENERATION=true
Environment=ENABLE_OBLIGATIONS_EXTRACTION=true
Environment=MCP_ONTOLOGY_INTEGRATION=true
Environment=USE_DB_VECTOR_SEARCH=true
Environment=LOG_LEVEL=INFO
Environment=DEBUG=false
ExecStart=/opt/proethica/venv/bin/python run.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=proethica-web

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/proethica/uploads /opt/proethica/logs /opt/proethica/backups
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=proethica.target