# server_config/mcp-server.service
[Unit]
Description=ProEthica MCP Server
After=network.target proethica-postgres.service
Requires=proethica-postgres.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/proethica
Environment="PATH=/var/www/proethica/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/var/www/proethica/venv/bin/python3 /var/www/proethica/mcp/run_enhanced_mcp_server.py
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
